# config/your_project.conf
[program:fastapi-app]
command=/platform/gunicorn_start
user=fastapi-user
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/platform/logs/gunicorn-error.log